# client/filmpinf.service
Description=filmpin frontend 
After=network.target multi-user.target

[Service]
User=mrps

WorkingDirectory=/home/mrps/prod/filmpin/client
ExecStart=node index.cjs
Environment=PORT=5175
StandardOutput=journal
StandardError=journal
SyslogIdentifier=filmpinf
Restart=Never

[Install]
WantedBy=multi-user.target
